# debian/agente-ia-defensiva-web.service
[Unit]
Description=Agente IA Defensiva Web Dashboard
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
WorkingDirectory=/opt/agente-ia-defensiva
Environment=PATH=/opt/agente-ia-defensiva/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/agente-ia-defensiva/.venv/bin/python /opt/agente-ia-defensiva/web/manage.py runserver 0.0.0.0:8000
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
