# debian/agente-ia-defensiva-ia.service
[Unit]
Description=Agente IA Defensiva Inference Engine
After=network.target postgresql.service agente-ia-defensiva-agent.service
Requires=postgresql.service
Wants=agente-ia-defensiva-agent.service

[Service]
Type=simple
WorkingDirectory=/opt/agente-ia-defensiva/src
Environment=PATH=/opt/agente-ia-defensiva/.venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/opt/agente-ia-defensiva/.venv/bin/python /opt/agente-ia-defensiva/src/ia_motor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
